# misc/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS Service

[Service]
ExecStart=/usr/bin/cloudflare-ddns -c /etc/cloudflare-ddns.yaml
