# misc/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS Service

[Service]
ExecStart=/usr/bin/cloudflare-ddns-service -c /etc/cloudflare-ddns.yaml
