# debian/cf-ddns.service
[Unit]
Description=Cloudflare DDNS Service

[Service]
ExecStart=/usr/bin/cf-ddns -c /etc/cf-ddns.yaml

[Install]
WantedBy=multi-user.target